# install/com.diachron.logrotate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.diachron.logrotate</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>
            LOG_DIR="$HOME/.diachron/logs"
            MAX_SIZE=$((10 * 1024 * 1024))  # 10MB
            MAX_FILES=7

            rotate_log() {
                local log="$1"
                if [ -f "$log" ]; then
                    local size=$(stat -f%z "$log" 2>/dev/null || echo 0)
                    if [ "$size" -gt "$MAX_SIZE" ]; then
                        for i in $(seq $((MAX_FILES - 1)) -1 1); do
                            [ -f "${log}.$i.gz" ] &amp;&amp; mv "${log}.$i.gz" "${log}.$((i + 1)).gz"
                        done
                        [ -f "${log}.1.gz" ] &amp;&amp; rm -f "${log}.$MAX_FILES.gz"
                        cp "$log" "${log}.1" &amp;&amp; gzip "${log}.1" &amp;&amp; : > "$log"
                    fi
                fi
            }

            rotate_log "$LOG_DIR/daemon.log"
            rotate_log "$LOG_DIR/daemon.err"
        </string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/dev/null</string>

    <key>StandardErrorPath</key>
    <string>/dev/null</string>
</dict>
</plist>
