# install/com.diachron.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.diachron.daemon</string>

    <key>ProgramArguments</key>
    <array>
        <string>$HOME/.claude/skills/diachron/rust/target/release/diachrond</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>StandardOutPath</key>
    <string>$HOME/.diachron/logs/daemon.log</string>

    <key>StandardErrorPath</key>
    <string>$HOME/.diachron/logs/daemon.err</string>

    <key>WorkingDirectory</key>
    <string>$HOME/.diachron</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>$HOME</string>
        <key>RUST_LOG</key>
        <string>info</string>
    </dict>

    <key>ProcessType</key>
    <string>Background</string>

    <key>ThrottleInterval</key>
    <integer>5</integer>
</dict>
</plist>
