# install/diachron-logrotate.timer
[Unit]
Description=Diachron Log Rotation Timer

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=3600

[Install]
WantedBy=timers.target
